# files/DO-NOT-USE-ME/systemd/user/wslg-weston.target
[Unit]
Description=Weston RDP (wslg)

BindsTo=wslg-weston.target
Before=wslg-weston.target
